# configs/autostart/first-boot.desktop
[Desktop Entry]
Name=first-boot
Comment=Run first-boot in terminal
Exec=~/.config/autostart/first-boot.sh
Terminal=true
Type=Application